# configs/pvpgn/systemd/pvpgn-pro.service
# /etc/systemd/system/pvpgn-pro.service
[Unit]
Description=PvPGN-PRO Battle.net Server
After=network.target

[Service]
# Run under a dedicated non-root user
User=root
Group=root

ProtectSystem=strict
ReadWritePaths=/usr/local/pvpgn/var

# Directory where PvPGN-PRO is installed
WorkingDirectory=/usr/local/pvpgn

# Path to binary + config file
ExecStart=/usr/local/pvpgn/sbin/bnetd -f -c /usr/local/pvpgn/etc/pvpgn/bnetd.conf

# Reload config without full restart (SIGHUP)
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target

